# files/system/etc/udev/rules.d/60-ac.rules
# This file was generated by IzZy to run a script when on ac

# Switch to balanced power profile
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/sbin/tuned-adm profile balanced"

# Start DisplayLink service
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/bin/systemctl start displaylink.service"