# files/system/etc/udev/rules.d/60-ac.rules
# This file was generated by IzZy to run a script when on ac

# Adjust Intel Power Limits Dell XPS 13 9310
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/bin/sh -c 'echo 20000000 > /sys/class/powercap/intel-rapl/intel-rapl:0/constraint_0_power_limit_uw'"
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/bin/sh -c 'echo 31000000 > /sys/class/powercap/intel-rapl/intel-rapl:0/constraint_1_power_limit_uw'"

# Switch to balanced power profile
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/sbin/tuned-adm profile balanced"

# Start DisplayLink service
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/bin/systemctl start displaylink.service"