# files/system/etc/udev/rules.d/60-battery.rules
# This file was generated by IzZy to run a script when on battery

# Stop DisplayLink service
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="0",RUN+="/bin/systemctl stop displaylink.service"

 Switch to battery power profile
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="0",RUN+="/usr/bin/system76-power profile battery"